# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autorestart=true

[program:php-fpm]
command=/usr/sbin/php-fpm7.4 -F
autorestart=true
